# docker_build_files/ope-smc/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/web2py.ini

[program:nginx-app]
command = /usr/sbin/nginx -g "daemon off;"