# ope-codecombat/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
loglevel=debug

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:coco]
user=root
directory=/home/coco
command=/home/coco/start.sh
stdout_logfile=/home/coco/codecombat/mongo/coco_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
priority=1000

[program:mongodb]
user=root
directory=/home/coco/codecombat
command=/home/coco/codecombat/bin/coco-mongodb
stdout_logfile=/home/coco/codecombat/mongo/mongodb_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
priority=100
