# docker_build_files/ope-fog/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
loglevel=debug

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D
autostart=false
autorestart=true


[program:apache]
command=apache2ctl -c "Error Log /dev/stdout" - DFOREGROUND
autostart=true
autorestart=true
startretries=3
startsecs=1
stdout_logfile=/var/log/apache_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
user=root
killasgroup=true
stopasgroup=true
priority=1000

[program:dhcpserver]
command=/usr/sbin/dhcpd -f
autostart=false
autorestart=true
startretries=3
startsecs=1
stdout_logfile=/var/log/dhcpd_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
user=root
killasgroup=true
stopasgroup=true
priority=900

[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysql.pid /usr/sbin/mysqld
autostart=true
autorestart=true
startretries=3
startsecs=1
stdout_logfile=/var/log/mysql_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
user=root
killasgroup=true
stopasgroup=true
priority=200


nfs-common
nfs-kernel-server
php5-fpm
tfptd-hpa
vsftp


[program:apache]

[program:mysql]


# REMOVE THIS - EXAMPLES

[program:canvas_web]
user=canvasuser
directory=/opt/canvas/canvas-lms
command=%(ENV_GEM_HOME)s/bin/bundle exec rails server -b 0.0.0.0
stdout_logfile=/opt/canvas/canvas-lms/log/server_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true

[program:canvas_worker]
user=canvasuser
directory=/opt/canvas/canvas-lms
command=%(ENV_GEM_HOME)s/bin/bundle exec script/delayed_job run
stdout_logfile=/opt/canvas/canvas-lms/log/delayed_job_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true


[program:mongodb]
user=root
directory=/home/coco/codecombat
command=/home/coco/codecombat/bin/coco-mongodb
stdout_logfile=/home/coco/codecombat/mongo/mongodb_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
priority=100


[program:coco]
user=root
directory=/home/coco
command=/home/coco/start.sh
stdout_logfile=/home/coco/codecombat/mongo/coco_stdout
stdout_logfile_maxbytes=50MB
redirect_stderr=true
priority=1000

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/web2py.ini

[program:nginx-app]
command = /usr/sbin/nginx -g "daemon off;"
