# docker_build_files/ope-canvas/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
username=docker
#loglevel=debug

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:canvas_web]
user=docker
directory=/usr/src/app
command=%(ENV_GEM_HOME)s/bin/bundle exec rails server -b 0.0.0.0
#stdout_logfile=/usr/src/app/log/server_stdout
#stdout_logfile_maxbytes=50MB
#redirect_stderr=true

[program:canvas_worker]
user=docker
directory=/usr/src/app
command=%(ENV_GEM_HOME)s/bin/bundle exec script/delayed_job run
#stdout_logfile=/usr/src/app/log/delayed_job_stdout
#stdout_logfile_maxbytes=50MB
#redirect_stderr=true
