# ope-canvas/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
loglevel=debug

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:canvas_web]
user=canvasuser
directory=/opt/canvas/canvas-lms
command=%(ENV_GEM_HOME)s/bin/bundle exec rails server
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:canvas_worker]
user=canvasuser
directory=/opt/canvas/canvas-lms
command=%(ENV_GEM_HOME)s/bin/bundle exec script/delayed_job run